# net/netatalk3/files/netatalk.in
#!/bin/sh
#
# $FreeBSD: net/netatalk3/files/netatalk.in 305024 2012-09-29 05:18:56Z marcus $
#
# PROVIDE: netatalk
# REQUIRE: DAEMON %%ZEROCONF%%
# KEYWORD: shutdown
#
# AFP fileserver for Mac clients.  Add the following to /etc/rc.conf to
# enable:
#
# netatalk_enable="YES"
#

netatalk_enable=${netatalk_enable-"NO"}

. /etc/rc.subr

name=netatalk
rcvar=netatalk_enable

load_rc_config ${name}

command=%%PREFIX%%/sbin/${name}

run_rc_command "$1"
