# net/netatalk3/files/netatalk.in
#!/bin/sh
#
# $FreeBSD: head/net/netatalk3/files/netatalk.in 340872 2014-01-24 00:14:07Z mat $
#
# PROVIDE: netatalk
# REQUIRE: DAEMON %%ZEROCONF%%
# KEYWORD: shutdown
#
# AFP fileserver for Mac clients.  Add the following to /etc/rc.conf to
# enable:
#
# netatalk_enable="YES"
#

netatalk_enable=${netatalk_enable-"NO"}

. /etc/rc.subr

name=netatalk
rcvar=netatalk_enable

load_rc_config ${name}

command=%%PREFIX%%/sbin/${name}

run_rc_command "$1"
