# comms/hylafax/files/hylafax.in
#!/bin/sh

# $FreeBSD: ports/comms/hylafax/files/hylafax.in,v 1.2 2012/11/17 05:54:59 svnexp Exp $
#
# PROVIDE: hylafax
# REQUIRE: LOGIN
# KEYWORD: shutdown

# be sure to initialize the following:
#   server: %%PREFIX%%/sbin/faxsetup
#   modems: %%PREFIX%%/sbin/faxaddmodem
#   incoming calls: /etc/ttys

hylafax_enable=${hylafax_enable-"NO"}

. /etc/rc.subr

name="hylafax"
rcvar="hylafax_enable"

start_cmd="%%PREFIX%%/sbin/hylafax start"
stop_cmd="%%PREFIX%%/sbin/hylafax stop"
procname="%%PREFIX%%/sbin/faxq"

load_rc_config $name
run_rc_command "$1"
