# mail/mailutils/files/pop3d.in
#!/bin/sh
#
# $FreeBSD: head/mail/mailutils/files/pop3d.in 429742 2016-12-28 12:01:54Z danfe $
#
# PROVIDE: pop3d
# REQUIRE: DAEMON
# BEFORE: mail
# KEYWORD: shutdown
#

. /etc/rc.subr

name="pop3d"
rcvar=pop3d_enable
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config ${name}

pop3d_enable=${pop3d_enable:-"NO"}
pop3d_flags=${pop3d_flags:-""}

run_rc_command "$1"
