# www/glassfish/files/glassfish.in
#!/bin/sh

# $FreeBSD: head/www/glassfish/files/glassfish.in 342332 2014-02-02 16:28:54Z swills $
#
# PROVIDE: glassfish
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# glassfish_enable (bool):	    Set to NO by default.
# 		  		    Set it to YES to enable doormand.

. /etc/rc.subr

name=glassfish
rcvar=glassfish_enable

load_rc_config $name

: ${glassfish_enable:=NO}

export PATH=${PATH}:%%LOCALBASE%%/bin

start_cmd="%%PREFIX%%/glassfish-%%GLASSFISHVERSION%%/bin/asadmin start-domain"
stop_cmd="%%PREFIX%%/glassfish-%%GLASSFISHVERSION%%/bin/asadmin stop-domain"

run_rc_command "$1"
