# net-p2p/bitcoin/files/bitcoind.in
#!/bin/sh
#
# $FreeBSD: head/net-p2p/bitcoin/files/bitcoind.in 414405 2016-05-01 09:50:43Z tijl $
#
# PROVIDE: bitcoind
# REQUIRE: LOGIN

. /etc/rc.subr

name="bitcoind"
rcvar=bitcoind_enable
command="%%PREFIX%%/bin/bitcoind"

load_rc_config $name
run_rc_command "$1"
