# news/nzbget/files/nzbget.in
#!/bin/sh
#
# $FreeBSD: head/news/nzbget/files/nzbget.in 434648 2017-02-22 22:24:16Z feld $
#
# PROVIDE: nzbget
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# nzbget_enable (bool):	Set to NO by default.
#			Set it to YES to enable it.

. /etc/rc.subr

name=nzbget
rcvar=nzbget_enable

load_rc_config ${name}

: ${nzbget_enable:=NO}

command=%%PREFIX%%/bin/nzbget
command_args="-D"

run_rc_command "$1"
