# audio/sndio/files/sndiod.in
#!/bin/sh
#
# $FreeBSD: head/audio/sndio/files/sndiod.in 476019 2018-07-31 19:40:24Z tobik $
#
# PROVIDE: sndiod
# REQUIRE: NETWORKING sysctl
# BEFORE:  DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=sndiod
rcvar=sndiod_enable

load_rc_config $name

: ${sndiod_enable="NO"}
: ${sndiod_flags="-c 0:7 -j off -s default -m mon -s monitor"}

command="%%PREFIX%%/bin/sndiod"

run_rc_command "$1"
