# security/clamsmtp/pkg-plist
@comment $FreeBSD: security/clamsmtp/pkg-plist 300897 2012-07-14 14:29:18Z beat $
@unexec [ ! -f %%CLAMAV_PID_DIR%%/clamsmtpd.pid ] || %%PREFIX%%/etc/rc.d/clamsmtpd.sh stop || true
sbin/clamsmtpd
etc/clamsmtpd.conf-sample
%%EXAMPLESDIR%%/virus_action.sh
@dirrm %%EXAMPLESDIR%%
