# net-mgmt/bandwidthd/pkg-plist
@comment $FreeBSD: net-mgmt/bandwidthd/pkg-plist 300897 2012-07-14 14:29:18Z beat $
bandwidthd/etc/bandwidthd.conf-dist
bandwidthd/bandwidthd
bandwidthd/htdocs/legend.gif
bandwidthd/htdocs/logo.gif
etc/rc.d/bandwidthd.sh.sample
@dirrmtry bandwidthd/etc
@dirrmtry bandwidthd/htdocs
@unexec rmdir %D/bandwidthd 2> /dev/null || echo "If you do not plan on reinstalling bandwidthd, you can safely remove %D/bandwidthd"
