# sysutils/cpupdate/files/cpupdate.in
#!/bin/sh

# $FreeBSD: head/sysutils/cpupdate/files/cpupdate.in 465350 2018-03-23 12:19:09Z eugen $
#
# PROVIDE: cpupdate
# REQUIRE: FILESYSTEMS kldxref
# BEFORE:  netif
# KEYWORD: nojail

# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# cpupdate_enable (bool):		Set to NO by default.
#					Set it to YES to enable cpupdate.
# cpupdate_flags (string):		Command line flags for cpupdate.
#					Default is "-w -u".

. /etc/rc.subr

name=cpupdate
rcvar=cpupdate_enable

load_rc_config $name

: ${cpupdate_enable:="NO"}
: ${cpupdate_flags="-w -u"}
command=%%PREFIX%%/sbin/${name}

run_rc_command "$1"
