# www/pound/files/pound.in
#!/bin/sh
#
# $FreeBSD: head/www/pound/files/pound.in 340872 2014-01-24 00:14:07Z mat $
#

# PROVIDE: %%NAME%%
# REQUIRE: NETWORKING SERVERS
# BEFORE:  DAEMON
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable `%%NAME%%':
#
#%%NAME%%_enable="YES"
#
# and be sure to configure %%PREFIX%%/etc/%%NAME%%.conf

. /etc/rc.subr

name="%%NAME%%"
rcvar=%%NAME%%_enable

command="%%PREFIX%%/sbin/$name"
pidfile="/var/run/$name.pid"
required_files="%%PREFIX%%/etc/$name.cfg"

# read configuration and set defaults
load_rc_config $name
: ${%%NAME%%_enable="NO"}

run_rc_command "$1"
