# sysutils/p5-MogileFS-Server/files/mogstored.in
#!/bin/sh

# $FreeBSD: sysutils/p5-MogileFS-Server/files/mogstored.in 310494 2013-01-16 15:07:59Z tobez $
#
# PROVIDE: mogstored
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown

. /etc/rc.subr

name=mogstored
rcvar=mogstored_enable

load_rc_config mogstored

: ${mogstored_enable:="NO"}

required_files=%%PREFIX%%/etc/mogstored.conf

command="%%PREFIX%%/bin/${name}"
command_args="--config=%%PREFIX%%/etc/mogstored.conf -d"

run_rc_command "$1"
