# databases/redis-devel/files/redis.in
#!/bin/sh
#
# $FreeBSD: head/databases/redis-devel/files/redis.in 369898 2014-10-03 13:27:09Z gahr $
#

# PROVIDE: redis
# REQUIRE: LOGIN
# BEFORE:  securelevel
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable `redis':
#
#redis_enable="YES"
#

. /etc/rc.subr

name="redis"
rcvar="${name}_enable"

extra_commands="reload"

command="%%PREFIX%%/bin/redis-server"
config_file="%%PREFIX%%/etc/$name.conf"
command_args="${config_file}"
pidfile="%%REDIS_RUNDIR%%/$name.pid"
required_files="${config_file}"

# read configuration and set defaults
load_rc_config "$name"
: ${redis_enable="NO"}
: ${redis_user="%%REDIS_USER%%"}

run_rc_command "$1"
