# dns/openmdns/files/mdnsd.in
#!/bin/sh
#
# $FreeBSD: head/dns/openmdns/files/mdnsd.in 424670 2016-10-26 04:33:37Z jbeich $
#
# PROVIDE: mdnsd
# REQUIRE: NETWORKING
# BEFORE:  DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="mdnsd"
desc="%%COMMENT%%"
rcvar="${name}_enable"
command="%%PREFIX%%/sbin/${name}"
command_args="em0"
extra_commands="reload"

: ${mdnsd_enable:="NO"}

load_rc_config $name
run_rc_command "$1"
