# net-mgmt/xymon-server/files/xymon-server.in
#!/bin/sh
#
# $FreeBSD: head/net-mgmt/xymon-server/files/xymon-server.in 346221 2014-02-26 18:09:20Z feld $
#

# PROVIDE: xymon_server
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=xymon_server
rcvar=xymon_server_enable

load_rc_config "$name"
: ${xymon_server_enable:=NO}
: ${xymon_server_user:=%%XYMONUSER%%}

command="%%WWWDIR%%/server/xymon.sh"
command_args=start
start_precmd=xymon_prestart
pidfile=/var/log/xymon/xymonlaunch.pid
procname="%%WWWDIR%%/server/bin/xymonlaunch"

xymon_prestart()
{
	if [ ! -d /var/log/xymon ]; then
		install -d -o xymon /var/log/xymon
	fi
}

run_rc_command "$1"
