# sysutils/usermin/files/usermin.in
#!/bin/sh
#
# $FreeBSD: sysutils/usermin/files/usermin.in 302141 2012-08-05 23:19:36Z dougb $
#
# PROVIDE: usermin
# REQUIRE: LOGIN
#
# Add the following line to /etc/rc.conf to enable usermin:
#
# usermin_enable="YES"
#

usermin_enable=${usermin_enable:-"NO"}

. /etc/rc.subr

name=usermin
rcvar=usermin_enable

procname=%%PERL5%%
pidfile=/var/log/usermin/miniserv.pid
required_dirs=%%PREFIX%%/etc/usermin
command=%%PREFIX%%/etc/usermin/start

load_rc_config ${name}
run_rc_command "$1"
