# sysutils/anacron/files/anacron.in
#!/bin/sh
#
# $FreeBSD: head/sysutils/anacron/files/anacron.in 340872 2014-01-24 00:14:07Z mat $
#

# PROVIDE: anacron
# REQUIRE: LOGIN

#
# Add the following line to /etc/rc.conf to enable anacron:
#
# anacron_enable="YES"
#

. /etc/rc.subr

name="anacron"
rcvar=anacron_enable
command="%%PREFIX%%/sbin/anacron"

load_rc_config $name
: ${anacron_enable="NO"}

run_rc_command "$1"
