# sysutils/volman/files/volmand.in
#!/bin/sh

# $FreeBSD: head/sysutils/volman/files/volmand.in 340872 2014-01-24 00:14:07Z mat $
#
# PROVIDE: volmand
# REQUIRE: LOGIN
# KEYWORD: shutdown nojail

. /etc/rc.subr

name="volmand"
rcvar=volmand_enable

load_rc_config $name

: ${volmand_enable:="NO"}

command="%%PREFIX%%/sbin/volmand"
pidfile="/var/run/volmand.pid"

run_rc_command "$1"
