# databases/tokyotyrant/files/ttserver.in
#!/bin/sh
#
# $FreeBSD: head/databases/tokyotyrant/files/ttserver.in 340872 2014-01-24 00:14:07Z mat $
#
#
# PROVIDE: ttserver
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="ttserver"
rcvar=ttserver_enable

: ${ttserver_enable:="NO"}

command="%%PREFIX%%/bin/${name}"
pidfile="/var/run/ttserver.pid"

load_rc_config $name
: ${ttserver_flags="-dmn -pid ${pidfile}"}

run_rc_command "$1"
