# net/chrony/files/chronyd.in
#!/bin/sh
#
# $FreeBSD: head/net/chrony/files/chronyd.in 454465 2017-11-18 22:41:34Z rodrigo $
#
# PROVIDE: chronyd
# REQUIRE: DAEMON
#

. /etc/rc.subr

name=chronyd
rcvar=chronyd_enable
command=%%PREFIX%%/sbin/${name}
rc_flags="${chronyd_flags}"

load_rc_config ${name}

: ${chronyd_enable="NO"}

run_rc_command "$1"
