# www/dhttpd/files/dhttpd.in
#!/bin/sh
#
# $FreeBSD: head/www/dhttpd/files/dhttpd.in 447117 2017-08-02 16:09:43Z lme $
#
# PROVIDE: dhttpd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable dhttpd:
#
# dhttpd_enable (bool):	Set it to "YES" to enable dhttpd
#

. /etc/rc.subr

name="dhttpd"
desc="Secure and efficient personal HTTP server"
rcvar=dhttpd_enable

load_rc_config $name

: ${dhttpd_enable="NO"}

command=%%PREFIX%%/sbin/${name}

run_rc_command "$1"
